# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock  ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
nodaemon=true
loglevel=debug
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket



[program:bk-cmdb]
directory=/data/abc/src/configcenter/src/bin/build/18.06.22/
command=/bin/bash start.sh
pidfile=/tmp//bk-cmdb.pid
autostart=true

